# systemd-user/backup.service
[Unit]
Description=Backup projects

[Service]
Type=oneshot
WorkingDirectory=%h
ExecStart=%h/bin/backup.sh
